# buildroot-external/board/weatherstation/rootfs_overlay/var/spool/cron/crontabs/root
*/5 * * * * /bin/python /root/log-data.py
